# profile.d/z_swift.sh
export PATH=/usr/local/opt/swift/bin:/Library/Developer/Toolchains/swift-latest.xctoolchain/usr/bin:$PATH
